# v0.4.3/vps/box/nfdos/servidor/nfdos/core/boot.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/no-boot.sh

[Install]
WantedBy=default.target